# supervisord.conf
; Supervisor config for the BamiHost FastAPI backend (pure-Python process manager).
; Replaces the old PM2 ecosystem.config.js — keeps the server always running with auto-restart.
;
; Usage (from the project root):
;   fastapi_app/venv/bin/supervisord -c supervisord.conf      # start manager + server
;   fastapi_app/venv/bin/supervisorctl -c supervisord.conf status
;   fastapi_app/venv/bin/supervisorctl -c supervisord.conf restart bami-backend
;   fastapi_app/venv/bin/supervisorctl -c supervisord.conf stop bami-backend

[supervisord]
logfile=%(here)s/logs/supervisord.log
pidfile=%(here)s/logs/supervisord.pid
childlogdir=%(here)s/logs
nodaemon=false

[unix_http_server]
file=%(here)s/logs/supervisor.sock

[supervisorctl]
serverurl=unix://%(here)s/logs/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:bami-backend]
command=%(here)s/fastapi_app/venv/bin/uvicorn main:app --host 0.0.0.0 --port 4000
directory=%(here)s/fastapi_app
autostart=true
autorestart=true
startretries=20
startsecs=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
environment=PYTHONPATH="%(here)s/fastapi_app"
stdout_logfile=%(here)s/logs/out.log
stderr_logfile=%(here)s/logs/error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
